# deploy/systemd/horizon-daily-publish.service
[Unit]
Description=Run Horizon daily digest and publish GitHub Pages
Wants=network-online.target docker.service
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=damon
Group=damon
SupplementaryGroups=docker
WorkingDirectory=/opt/projects/Horizon
Environment=HOME=/home/damon
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=TZ=Pacific/Auckland
Environment=HORIZON_HOURS=24
ExecStart=/opt/projects/Horizon/scripts/run-and-publish.sh
TimeoutStartSec=12h
