# deploy/systemd/user/horizon-daily-publish.timer
[Unit]
Description=Run Horizon daily digest and publish GitHub Pages at 08:00 NZ time

[Timer]
OnCalendar=*-*-* 08:00:00 Pacific/Auckland
Persistent=true
Unit=horizon-daily-publish.service

[Install]
WantedBy=timers.target
